# .config/systemd/user/cmus-update-cache.service
[Unit]
Description = Update cmus cache

[Service]
Type=oneshot
ExecStart = /usr/bin/cmus-remote --raw 'update-cache'
