# .config/systemd/user/tldr.service
[Unit]
Description=Update tldr cache

[Service]
Type=oneshot
ExecStart=/usr/bin/tldr --update
